# deploy/tori.service
[Unit]
Description=Tori Server Monitoring Agent
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=tori
Group=tori
ExecStart=/usr/local/bin/tori agent --config /etc/tori/config.toml
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5

RuntimeDirectory=tori
RuntimeDirectoryMode=0750

# Security hardening
ProtectSystem=strict
ReadWritePaths=/var/lib/tori
ProtectHome=true
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
